# config/systemd/user/urxvtd.service
[Unit]
Description=Urxvt daemon

[Service]
Environment=SHELL=/bin/zsh
ExecStart=/usr/bin/urxvtd

[Install]
WantedBy=mystuff.target
